# Ass04/ldamulel-usb.rules
SUBSYSTEM=="input", SUBSYSTEMS=="usb", ACTION=="add", RUN+="sbin/insmod /root/ldamulel/little_p/Ass04/main.ko"
